# rootfs/etc/services.d/sleep/run
#!/command/execlineb
/bin/sleep infinity
